# celery.conf
[program:celery]
command=celery -A app.game worker -b redis:// --loglevel=debug
directory=/app
