# celery.conf
[program:celery]
command=celery -A app.game worker -E -b redis:// --loglevel=info
directory=/app
